# systemd/topdown-watch.service
[Unit]
Description=Quant Fund — top-down brief on-change watcher (push on bias flip)
Documentation=file:///home/jer/quant-fund/docs/SYSTEM_BASELINE.md
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/jer/quant-fund
EnvironmentFile=/home/jer/quant-fund/config/.env
Environment=PYTHONPATH=/home/jer/quant-fund
ExecStart=/home/jer/quant-fund/venv/bin/python /home/jer/quant-fund/scripts/topdown_push.py watch --interval 15
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
Nice=10

[Install]
WantedBy=default.target
